# airootfs/etc/xdg/autostart/cnchi.desktop
[Desktop Entry]
Type=Application
Name=Reborn Installer
Exec=/usr/bin/antergos-welcome
Icon=/usr/share/applications/cnchi.png
Terminal=True
