# airootfs/etc/xdg/autostart/internet.desktop
[Desktop Entry]
Type=Application
Name=Refresh Pacman Databases
Exec=bash /usr/bin/internet.sh
