# airootfs/etc/xdg/autostart/rank.desktop
[Desktop Entry]
Type=Application
Name=Ranking Mirrors Popup
Exec=bash /usr/bin/cnchi-rank.sh
